# webapp/analytics/goaccess-analytics.timer
# Runs the analytics refresh every 10 minutes (and shortly after boot).
# See goaccess-analytics.service.

[Unit]
Description=Periodic GoAccess analytics refresh

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
